# centos/conf/supervisor.d/httpd.conf
[program:apache2]
directory=/var/www/owlh
command=/usr/sbin/apachectl -DFOREGROUND
autostart=true
autorestart=true
redirect_stderr=true
process_name=httpd
