# centos/conf/supervisor.d/owlhmaster.conf
[program:owlhmaster]
directory=/usr/local/owlh/src/owlhmaster
command=/usr/local/owlh/src/owlhmaster/owlhmaster
autostart=true
autorestart=true
redirect_stderr=true

