# centos/conf/supervisor.d/owlhnode.conf
[program:owlhnode]
directory=/usr/local/owlh/src/owlhnode
command=/usr/local/owlh/src/owlhnode/owlhnode
autostart=true
autorestart=true
redirect_stderr=true

